# etcfs/service/apache2/run
#!/usr/bin/env bash
source /etc/apache2/envvars
source /etc/envvars
exec   /usr/sbin/apache2 -DFOREGROUND
